# debian/spampd.service
[Unit]
Description=Spam Proxy Daemon
After=network.target local-fs.service

[Service]
EnvironmentFile=/etc/default/spampd
ExecStart=/usr/sbin/spampd -u $USERID -g $GRPID $SPAMPD_ARGS $ADDOPTS --pid=/var/run/spampd/spampd.pid
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
PIDFile=/var/run/spampd/spampd.pid

[Install]
WantedBy=multi-user.target
