# debian/spampd.service
[Unit]
Description=Spam Proxy Daemon
After=network.target local-fs.service

[Service]
ExecStart=/usr/sbin/spampd --setsid
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
PIDFile=/var/run/spampd/spampd.pid

[Install]
WantedBy=multi-user.target
